# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS updater for sh.ketencek.com.au
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=salih
Group=salih
WorkingDirectory=/home/salih/projects/cloudflare-ddns
EnvironmentFile=/home/salih/projects/cloudflare-ddns/.env
ExecStart=/home/salih/projects/cloudflare-ddns/venv/bin/python /home/salih/projects/cloudflare-ddns/ddns.py

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/salih/projects/cloudflare-ddns

[Install]
WantedBy=multi-user.target
